# pkgs/home-automount-nas/mnt-itsuki.automount
[Unit]
After=NetworkManager.service

[Automount]
Where=/mnt/itsuki
TimeoutIdleSec=0

[Install]
WantedBy=graphical.target