# run/robocadey.service.in
[Unit]
Description=Robo-Cadey

[Service]
ExecStartPre=@robocadey@/bin/robocadey --read-from=./robocadey.gob
ExecStart=@robocadey@/bin/robocadey
StateDirectory=xeserv.robocadey
CacheDirectory=xeserv.robocadey

[Install]
WantedBy=multi-user.target
